# .docker/root/etc/services.d/medusa/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/medusa/start.py \
	--nolaunch --datadir /config
